# runtime/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
chown=node:node

[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:kasmvnc]
command=/bin/bash -lc "/usr/local/bin/kasmvnc_start"
user=node
priority=10
autostart=false
autorestart=true
startretries=10
environment=HOME="/home/node"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:gost]
command=/bin/bash -lc "/usr/local/bin/gost_start"
priority=35
autostart=false
autorestart=true
startretries=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:state-dir-fix]
command=/bin/bash -lc 'd="${SELLERCLAW_DATA_DIR:-/data}"; mkdir -p "$d" /home/node/.openclaw /home/node/.cache && chown -R node:node "$d" /home/node/.openclaw /home/node/.cache'
priority=38
autostart=true
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:openclaw]
command=/bin/bash -lc "/usr/local/bin/openclaw_start gateway"
user=node
priority=40
autostart=true
autorestart=true
startretries=10
environment=DISPLAY=":1",HOME="/home/node"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
